# s6/gost/run
#!/bin/sh

CONFFILE=/config/gost.yml
while ! [ -e "$CONFFILE" ]; do
  sleep 5
done

cd /tmp
exec /app/gost "$CONFFILE"
