# s6/gost/run
#!/bin/sh

CONFFILE=/config/gost.yml
while ! [ -e "$CONFFILE" ]; do
  sleep 5
done

cd /tmp
rm -rf gost
mkdir -p gost
cd gost
exec /app/gost -C "$CONFFILE"
