# s6/corplink/run
#!/bin/sh

exec /app/corplink-rs /config/config.json
